# share/avst-app/upstart/upstart-debian.conf
# avst-app - Adaptavist Application - Services
#
# By default these are in: /opt/<instance name>

description "Avst-app Instances"

# Start once we have networking and filesystems
start on (local-filesystems and net-device-up IFACE!=lo)

# Stop when shutting down the system or entering single user mode
stop on runlevel [016]

# Auto-magically restart on exit
respawn

# Respawn the job up to 10 times within a 5 second period
respawn limit 10 5

# Try hard not to be killed in Out-Of-Memory situations
oom score -999

# Wait before shooting app (in seconds)
kill timeout 120

# Log to /var/log/upstart
console log

# UPSTART_JOB is set to the filename, which will be the Instance Name
exec avst-app $UPSTART_JOB start

# sleep for 5 seconds after stop, work around for fast restarts
post-stop exec sleep 5
